# proxmox/cron.daily
#!/bin/sh
/usr/sbin/smartctl -t long /dev/sda
/usr/sbin/smartctl -t long /dev/sdb
