# proxmox/cron.daily
#!/bin/sh
/usr/sbin/smartctl -t long /dev/sda >/dev/null
/usr/sbin/smartctl -t long /dev/sdb >/dev/null
