# proxmox/cron.weekly
#!/bin/sh
# trim all the mounted file systems which support it
/usr/sbin/fstrim --all || true
